# roles/prometheus/files/prometheus-node-exporter.service
[Service]
ExecStart=
ExecStart=/usr/bin/prometheus-node-exporter --collector.systemd $ARGS

